# systemd/rclone-openlist-root.service
[Unit]
Description=Rclone Mount OpenList WebDAV Root
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /srv/cloud
ExecStart=/usr/bin/rclone --config /etc/rclone/rclone.conf mount openlist:/ /srv/cloud \
  --allow-other \
  --umask 002 \
  --vfs-cache-mode full \
  --cache-dir /var/cache/rclone \
  --dir-cache-time 10m \
  --timeout 1h
ExecStop=/bin/sh -c 'fusermount3 -uz /srv/cloud || fusermount -uz /srv/cloud'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
